# ubuntu16/cloudboost.service
[Unit]
Description=Start CloudBoost Server

[Service]
ExecStart=/bin/bash startup.sh start
Type=oneshot

[Install]
WantedBy=multi-user.target
